# tools/com.beako.procfs.mount.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Per-user LaunchAgent: mounts procfs at ~/proc on login (gated; runs
  procfs-mount.sh which no-ops unless /var/db/procfs.enabled exists). Install to
  /Library/LaunchAgents/com.beako.procfs.mount.plist (root:wheel, 644), then:
      launchctl bootstrap gui/$(id -u) /Library/LaunchAgents/com.beako.procfs.mount.plist
  The system LaunchDaemon (com.beako.procfsd) handles symbol staging + kext load
  at boot, so by login the kext is loaded and the mount can proceed.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.beako.procfs.mount</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/sbin/procfs-mount.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/procfs-mount.log</string>
</dict>
</plist>
